# deploy/systemd/merdian-wsfeed-stop.service
[Unit]
Description=MERDIAN wsfeed stop (market close)

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl stop merdian-wsfeed.service
